# current/debian/passwd.conffiles
/etc/cron.daily/passwd
/etc/init.d/passwd
/etc/shells
