# debian/passwd.cron
#!/bin/sh
#
# cron.daily script to check integrity of the password and group files

test -f /usr/sbin/pwck || exit 0

pwck -q -r
grpck -r
